# crontab
# Onda 3 — crontab do runner (lido pelo supercronic no Fly.io). Formato cron padrão de 5 campos.
# Horários em UTC (TZ=UTC no fly.toml). Cada linha redireciona para um log em /app/logs.

# Fila: processa no máximo 1 job por minuto (lock + claim atômico em poll-agent-jobs.sh).
* * * * * bash /app/scripts/poll-agent-jobs.sh >> /app/logs/poll.log 2>&1

# Onda 9 — Modo autônomo: avança 1 tick de um watch ativo por minuto (lock + claim atômico).
* * * * * bash /app/scripts/poll-autonomous-watches.sh >> /app/logs/watch.log 2>&1

# Onda 12 — Valida os tokens Meta manuais das conexões vivas (1×/dia, SPEC §5.4). TS determinístico,
# sem claude -p. Antes da skill de tráfego para pegar um token revogado cedo. Precisa AD_TOKEN_ENC_KEY.
0 8 * * * cd /app && npx tsx scripts/onda12/validate-connections.ts >> /app/logs/cron.log 2>&1

# Skill de tráfego agendada — aceite da Onda 3: "cron dispara a skill da Onda 2". Diária às 09:00 UTC.
0 9 * * * bash /app/scripts/run-skill.sh create-traffic-cliente-exemplo-campaign >> /app/logs/cron.log 2>&1

# Onda 4 — Analytics read-only: extrai o funil + diagnóstico do período anterior. Diária às 10:00 UTC.
0 10 * * * bash /app/scripts/run-skill.sh funnel-analytics-cliente-exemplo-campaign >> /app/logs/cron.log 2>&1

# Onda 4 — Resumo diário: agrega as analyses do dia (depende da análise acima). Diária às 10:30 UTC.
30 10 * * * bash /app/scripts/run-skill.sh daily-summary-cliente-exemplo >> /app/logs/cron.log 2>&1
